# systemd/magpie-firewall.service
# magpie-firewall.service — boot-time network lockdown for the reviewer
# (Milestone 5, task_56ad; wraps scripts/setup-network.sh from M4-D).
#
# A oneshot that provisions the locked-down `magpie-net` docker bridge and the
# host iptables rules pinning it to the credential-injecting gateway (see
# scripts/setup-network.sh's header and PLAN.md §5). It is ORDERED BEFORE both
# the gateway and the orchestrator on purpose:
#   - the gateway binds its proxy plane to the bridge's host IP (172.31.99.1),
#     which does not exist until this unit has created br-magpie; and
#   - no review container should be launchable until egress is locked down.
# setup-network.sh is idempotent, so re-running on every boot (or a manual
# `systemctl restart magpie-firewall`) is safe and never duplicates a rule.
#
# Install: see scripts/install.sh (copies this into /etc/systemd/system).

[Unit]
Description=magpie network lockdown (magpie-net bridge + egress iptables)
Documentation=file:/opt/magpie/scripts/setup-network.sh
# Needs docker (it creates the bridge network) and a configured network stack.
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target
# The gateway and orchestrator must not start until the bridge + rules exist.
Before=magpie-gateway.service magpie.service

[Service]
Type=oneshot
# Keep the unit "active (exited)" after the script returns so the ordering
# above holds and `systemctl status` reflects that the lockdown is in place.
RemainAfterExit=yes
# Runs as root: setup-network.sh needs `docker network` + `iptables`. The
# script is `set -euo pipefail` and fails closed, so any error applying a rule
# leaves this unit failed and (via the Before= ordering) stops the gateway and
# orchestrator from coming up on a half-locked bridge.
ExecStart=/opt/magpie/scripts/setup-network.sh
# Intentionally NOT sandboxed (no ProtectSystem/RestrictNamespaces/etc.): this
# unit's whole job is to manipulate netfilter and the docker bridge, which
# needs broad privilege (CAP_NET_ADMIN, module autoload). Its blast radius is
# bounded by being a short-lived oneshot running a source-controlled, idempotent
# script — the long-running, attacker-adjacent processes (gateway, reviewer
# container) are the ones that get the tight sandboxes.

[Install]
WantedBy=multi-user.target
