# systemd/magpie-gateway.service
# magpie-gateway.service — the credential-injecting LLM gateway
# (Milestone 5, task_56ad; runs @magpie/gateway from packages/gateway, M4-A;
# per-job unix-socket proxy plane added M7-1, Design D — DISTRIBUTION.md §2).
#
# This is the ONLY process that holds the real OpenRouter key
# (MAGPIE_GATEWAY_OPENROUTER_KEY). It runs as its OWN unprivileged user,
# `magpie-gateway`, separate from the orchestrator's `magpie` user, so a
# compromise of the orchestrator's filesystem access never yields the provider
# key (see packages/gateway/README.md "Provisioning" and PLAN.md §5). It has
# NO docker access — the virtual-key store is in-memory by design (keys are
# per-job ephemeral; losing them on restart just fails any in-flight job
# closed), so there is no persistent state directory to protect. It DOES now
# own one writable-but-ephemeral path: the per-job proxy-plane socket
# directory (`RuntimeDirectory=` below), which holds only unix socket special
# files, never secrets or review content.
#
# Proxy-plane transport (M7-1): as of Design D there is no more magpie-net
# bridge / fixed proxy IP to bind. Each job gets its own unix domain socket
# under `RuntimeDirectory=magpie-gateway` (systemd creates and owns
# `/run/magpie-gateway`, mode 0700, `magpie-gateway:magpie-gateway`); this
# process creates the per-job `jobs/<id>/` subdirectories (mode 0711) and
# binds `jobs/<id>/gw.sock` (chmod 0666 after bind) itself — see
# DISTRIBUTION.md §2.6 point 2 for the full launch-ordering/permissions
# rationale, and packages/gateway/src/job-sockets.ts for the implementation.
# The docker daemon, running as root on the host, is what actually makes a
# given job's socket directory visible inside that job's `--network none`
# reviewer container: the orchestrator bind-mounts `jobs/<id>/` read-only at
# `/run/gw` in that one container (see packages/orchestrator/src/reviewer.ts).
# The mgmt plane (127.0.0.1:4100, mint/revoke) is unchanged — still
# loopback-only, still how the orchestrator talks to this process.
#
# Ordering: no longer depends on any firewall/network-lockdown unit (that
# apparatus — magpie-firewall.service, scripts/setup-network.sh — is deleted
# entirely under Design D: the reviewer has no network, so there is nothing
# left to provision). Still started BEFORE the orchestrator so a PR review can
# mint a key the moment magpie is up.
#
# Secrets come from EnvironmentFile=/etc/magpie-gateway/gateway.env (chmod 600,
# owned by magpie-gateway) — see scripts/install.sh. Required keys:
#   MAGPIE_GATEWAY_OPENROUTER_KEY  (the real provider key)
#   MAGPIE_GATEWAY_MASTER_KEY      (shared with the orchestrator; mgmt-plane auth)
# and optionally GATEWAY_DEFAULT_MODEL, GATEWAY_SOCKET_DIR, GATEWAY_MGMT_PORT.

[Unit]
Description=magpie LLM gateway (credential-injecting OpenRouter proxy)
Documentation=file:/opt/magpie/packages/gateway/README.md
After=network-online.target
Wants=network-online.target
Before=magpie.service

[Service]
Type=simple
User=magpie-gateway
Group=magpie-gateway
WorkingDirectory=/opt/magpie/packages/gateway
EnvironmentFile=/etc/magpie-gateway/gateway.env
# RuntimeDirectory: systemd creates /run/magpie-gateway (0700,
# magpie-gateway:magpie-gateway) before this unit starts and removes it on
# stop — exactly the top-level layout DISTRIBUTION.md §2.6 point 2 specifies.
# This process creates the `jobs/` subdir and each per-job `jobs/<id>/`
# beneath it itself (job-sockets.ts), at the narrower 0711 the design calls
# for. GATEWAY_SOCKET_DIR is set explicitly here for clarity even though it
# matches the gateway's own default (packages/gateway/src/config.ts's
# DEFAULT_SOCKET_DIR_ROOT) — an operator overriding RuntimeDirectory's name
# would otherwise silently desync the two.
RuntimeDirectory=magpie-gateway
Environment=GATEWAY_SOCKET_DIR=/run/magpie-gateway/jobs
ExecStart=/usr/bin/node /opt/magpie/packages/gateway/dist/index.js
Restart=on-failure
RestartSec=5s

# --- systemd hardening ------------------------------------------------------
# The gateway is host-side but sits directly behind the (fully-compromisable)
# review container on the data plane, so it gets a tight sandbox. It needs
# only: outbound TCP to OpenRouter, a listening TCP socket for the mgmt plane,
# a listening unix socket per job under RuntimeDirectory for the proxy plane,
# and read access to its own code + env file. Nothing else.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# No writable filesystem paths beyond RuntimeDirectory=magpie-gateway above
# (systemd grants that one exception automatically): the key store is
# in-memory, so there is nothing else to persist. ProtectSystem=strict makes
# the rest of the tree read-only. Drop ALL capabilities — the mgmt (4100) TCP
# port is >1024 and unix sockets need no capability at all, so not even
# CAP_NET_BIND_SERVICE is needed.
CapabilityBoundingSet=
AmbientCapabilities=
# INET(6) for the mgmt-plane listener + the outbound call to OpenRouter; UNIX
# for the per-job proxy-plane sockets under RuntimeDirectory (M7-1) and
# Node's internal use (e.g. DNS/threadpool signalling); NETLINK because
# glibc's resolver (getaddrinfo) opens an AF_NETLINK socket to enumerate
# interfaces for source-address/IPv6 selection — blocking it can make DNS
# lookups (e.g. openrouter.ai) slow or fail. No others.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
SystemCallArchitectures=native
# NOTE: deliberately NO MemoryDenyWriteExecute — V8's JIT needs W^X memory and
# would crash under it.

[Install]
WantedBy=multi-user.target
